# launcherunitfile.service
[Unit]
Description=Launcher

[Service]
ExecStart=bash /home/claude/a/launcher
Type=simple
